# power-profiles-cfg.service.in
[Unit]
Description=Power Profile Configuration
After=power-profiles-daemon.service

[Service]
Type=simple
ExecStart=@LIBEXECDIR@/power-profiles-cfg --init
Restart=on-failure
StateDirectory=power-profiles-cfg

[Install]
WantedBy=multi-user.target